# magazyn360-api/docs/supervisord/magazyn360.conf
[supervisord]
nodaemon=true
logfile=/var/log/magazyn360/supervisord.log
pidfile=/tmp/supervisord.pid

[program:magazyn360-gunicorn]
command=/bin/bash /var/www/magazyn360/app/magazyn360-api/bin/magazyn360_start_gunicorn.sh
directory=/var/www/magazyn360/app/magazyn360-api/
autostart=true
autorestart=true
stdout_logfile=/var/log/magazyn360/startup.log
stderr_logfile=/var/log/magazyn360/startup_err.log
